# deployment/systemd/reverb.service
[Unit]
Description=Laravel Reverb WebSocket Server - Hearts Connect Dating Platform
After=network.target mysql.service redis.service
Requires=network.target

[Service]
Type=simple
User=www-data
Group=www-data
Restart=always
RestartSec=5
StartLimitInterval=0

# Execution
ExecStart=/usr/bin/php /var/www/html/dating/artisan reverb:start
WorkingDirectory=/var/www/html/dating

# Logging
StandardOutput=append:/var/www/html/dating/storage/logs/reverb.log
StandardError=append:/var/www/html/dating/storage/logs/reverb-error.log
SyslogIdentifier=reverb

# Security hardening
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/html/dating/storage

# Environment
Environment="LANG=en_US.UTF-8"
Environment="LC_ALL=en_US.UTF-8"

[Install]
WantedBy=multi-user.target
